# container-files/config/gitlab-runner.conf
[program:gitlab-runner]
command = /bin/gitlab-runner run --user=gitlab-runner --working-directory=/home/gitlab-runner --syslog --config=/etc/gitlab-runner/config.toml
autorestart = false
priority = 999
stdout_events_enabled=true